# config/systemd/user/rclone@.service
[Unit]
Description=Rclone mount for %i
Documentation=http://rclone.org/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
# Check if mount point exists, if not create it
ExecStartPre=/usr/bin/mkdir -p %h/%i
# Mount command
ExecStart=/usr/bin/rclone mount \
    --vfs-cache-mode full \
    --vfs-cache-max-size 10G \
    %i: %h/%i
# Unmount command
ExecStop=/bin/fusermount -u %h/%i
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
